# Linux/ipod/ipod-iap2/test/iap2_smoketest/iap2_bluetooth.service
[Unit]
Description=Bluetooth service

[Service]
Type=dbus
BusName=org.bluez
ExecStart=/usr/lib/bluez5/bluetooth/bluetoothd -n

[Install]
WantedBy=multi-user.target
