# Linux/ipod/ipod-common/src/auth_plugin/auth_sub20_plugin/ipod-auth.target
[Unit]
Description=iPod authentication
After=
Requires=

