# Linux/multimedia/platform/audiomanager/ias-audio-configuration_brd2/src-gen/ias-audio-configuration_brd2-10.0.0/public/brd2/systemd/audio_daemon.path
[Path]
PathExists=/dev/shm/jack_default_0_0
Unit=audio_daemon.service

