# Linux/ipod/ipod-iap2/test/iap2_service_smoketest/iap2.socket
[Socket]
#SocketUser=system
#SocketGroup=system
SocketMode=0660
ListenSequentialPacket=/tmp/iAP2_Server

[Install]
WantedBy=sockets.target
