# Linux/ipod/ipod-common/src/auth_plugin/auth_i2c_plugin/ipod-auth.target
[Unit]
Description=iPod authentication
After=i2c-dev-setup.service
Requires=i2c-dev-setup.service

